# pi/webui/homelab-ui.service.template
[Unit]
Description=Homelab control plane (FastAPI + SPA)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__WEBUI__
ExecStart=__WEBUI__/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 80 --timeout-graceful-shutdown 2
Restart=on-failure
RestartSec=2
TimeoutStopSec=3
KillMode=mixed
User=root

[Install]
WantedBy=multi-user.target
