# pi/webui/homelab-poll.service.template
[Unit]
Description=Homelab device uptime poller
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-__REPO__/secrets.env
ExecStart=__WEBUI__/venv/bin/python __WEBUI__/poll.py
