# pi/webui/homelab-ui.service.template
[Unit]
Description=Homelab control plane web UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=__WEBUI__/venv/bin/python __WEBUI__/app.py
WorkingDirectory=__WEBUI__
Restart=on-failure
RestartSec=2
User=root

[Install]
WantedBy=multi-user.target
